# systemd/pkgbuilder.timer
[Unit]
Description=pkgbuilder - check packages update

[Timer]
OnCalendar=*-*-* 16:00:00

[Install]
WantedBy=timers.target