# systemd/pkgbuilder.service
[Unit]
Description=pkgbuilder

[Service]
Type=simple
ExecStart=/usr/bin/pkgbuilder checkupdateall --notify
